# stemcell_builder/stages/centos_bosh_agent/assets/bosh_agent.conf
# for runit - manage /usr/sbin/boshagent-start
start on runlevel [2345]
stop on shutdown
respawn
exec /usr/sbin/boshagent-start
